# etc/systemd/wickedd.service.in
[Unit]
Description=wicked network management service daemon
Requisite=dbus.service
Wants=wickedd-nanny.service wickedd-dhcp6.service wickedd-dhcp4.service wickedd-auto4.service
After=local-fs.target dbus.service isdn.service rdma.service SuSEfirewall2_init.service
Before=wickedd-nanny.service wicked.service network.target

[Service]
Type=notify
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/network/config
ExecStartPre=-/usr/bin/udevadm settle
ExecStart=@wicked_sbindir@/wickedd --systemd --foreground
StandardError=null
Restart=on-abort

[Install]
Also=wickedd-nanny.service
Also=wickedd-auto4.service
Also=wickedd-dhcp4.service
Also=wickedd-dhcp6.service
