# etc/systemd/wickedd-nanny.service.in
[Unit]
Description=wicked network nanny service
Requisite=dbus.service
After=local-fs.target dbus.service SuSEfirewall2_init.service wickedd.service
Before=wicked.service network.target
PartOf=wickedd.service

[Service]
Type=notify
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/network/config
ExecStart=@wicked_sbindir@/wickedd-nanny --systemd --foreground
StandardError=null
Restart=on-abort

[Install]
Alias=dbus-org.opensuse.Network.Nanny.service

