# etc/systemd/wickedd-dhcp6.service.in
[Unit]
Description=wicked DHCPv6 supplicant service
Requisite=dbus.service
After=local-fs.target dbus.service SuSEfirewall2_init.service
Before=wickedd.service wicked.service network.target
PartOf=wickedd.service

[Service]
Type=notify
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/network/config
ExecStart=@wicked_supplicantdir@/wickedd-dhcp6 --systemd --foreground
StandardError=null
Restart=on-abort

[Install]
Alias=dbus-org.opensuse.Network.DHCP6.service

