# etc/systemd/wicked@.service.in
#
# This is a template example intended to implement
# a multi-stage network start.
#
[Unit]
Description=wicked managed network in stage %i
#Wants=network.target
#Before=network.target
After=wickedd.service

[Service]
Type=oneshot
RemainAfterExit=yes
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/network/config
ExecStart=@wicked_sbindir@/wicked ifup --boot-stage %i -- all
ExecStop=@wicked_sbindir@/wicked ifdown --boot-stage %i -- all

[Install]
Also=wickedd.service
#WantedBy=some.target

