# etc/systemd/wickedd.service.in
[Unit]
Description=wicked network management service daemon
After=local-fs.target dbus.service isdn.service openibd.service
Before=wicked.service

[Service]
Type=dbus
BusName=org.opensuse.Network
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/network/config
ExecStart=@wicked_sbindir@/wickedd $WICKED_DEBUG_PARAM --foreground
StandardError=null

[Install]
WantedBy=wicked.service
Also=wickedd-nanny.service
Also=wickedd-auto4.service
Also=wickedd-dhcp4.service
Also=wickedd-dhcp6.service
