# etc/systemd/wicked.service.in
[Unit]
Description=wicked managed network interfaces
Wants=network.target
After=wickedd.service wickedd-nanny.service
Before=SuSEfirewall2.service network-online.target network.target multi-user.target shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
LimitCORE=infinity
EnvironmentFile=-/etc/sysconfig/network/config
ExecStart=@wicked_sbindir@/wicked --systemd ifup all
ExecStop=@wicked_sbindir@/wicked --systemd ifdown all
ExecReload=@wicked_sbindir@/wicked --systemd ifreload all

[Install]
WantedBy=multi-user.target
Alias=network.service
Also=wickedd.service

